# s6/syslogd/run
#!/usr/bin/env bash

[ "$DEBUG" == 'true' ] && set -x

exec syslogd -n -O /dev/stdout -S
