# s6/syslogd/run
#!/usr/bin/env bash

[ "$DEBUG" == 'true' ] && set -x

# Defaults
: ${LOGOUTPUT:='/dev/stdout'}

echo "Setting syslogd output to ${LOGOUTPUT}"

exec syslogd -n -O ${LOGOUTPUT} -S
