# etc/s6/syslogd/run
#!/usr/bin/env bash

set -e

[ "$DEBUG" == 'true' ] && set -x

echo "syslogd >> Setting syslogd output to /dev/stdout (Only used for OpenDKIM)"

exec syslogd -n -O "/dev/stdout" -S
