# s6/rsyslog/run
#!/usr/bin/env bash

[ "$DEBUG" == 'true' ] && set -x

exec /usr/sbin/rsyslogd -n
