# s6/catlog/finish
#!/bin/sh

# Shutdown everything and exit the process crashes or is stopped.
s6-svscanctl -t /etc/s6
