# s6/catlog/run
#!/usr/bin/env bash

[ "$DEBUG" == 'true' ] && set -x

# Wait for postfix to create FIFO
while [ ! -e "/dev/maillog" ]; do
    sleep 1
done

exec cat /dev/maillog
