# s6/opendkim/finish
#!/bin/sh

echo "FINISH $(pwd) WITH: ${@}"
