# s6/postfix_exporter/finish
#!/usr/bin/env bash

if [[ "${POSTFIX_EXPORTER_ENABLED}" == "true" ]]; then
  # Shutdown everything and exit the process crashes or is stopped.
  s6-svscanctl -t /etc/s6
fi
